# crates/vtt-mcp/vtt-mcp.service
[Unit]
Description=VTT-MCP Voice-to-Text Server
Documentation=man:vtt-mcp(1)
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=vtt-mcp
Group=vtt-mcp
WorkingDirectory=/var/lib/vtt-mcp
ExecStart=/usr/bin/vtt-mcp
Restart=on-failure
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vtt-mcp /var/log/vtt-mcp

# Audio access
DeviceAllow=/dev/snd rw

[Install]
WantedBy=multi-user.target
